# APT_Groups/Scattered_Spider_/T1053_-Scheduled_Task_Creation.kql
// T1053 - Scheduled Task Creation
// Detects suspicious task creation for persistence
DeviceProcessEvents
| where FileName == "schtasks.exe" and ProcessCommandLine contains "/create"
| project Timestamp, DeviceName, ProcessCommandLine
