# APT Groups/APT29 (Cozy Bear)/Techniques/T1053.005 - Scheduled Task.kql
// T1053.005 - Scheduled Task
// Detect creation of scheduled tasks via schtasks or PowerShell
DeviceProcessEvents
| where ProcessCommandLine contains_any ("schtasks.exe /create", "New-ScheduledTaskTrigger")
| project Timestamp, DeviceName, ProcessCommandLine, InitiatingProcessFileName
| order by Timestamp desc
